# iqoptiontraderbot.service
[Unit]
Description=IQ Option Trader Bot - Gunicorn daemon
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/iqoptiontraderbot
Environment="PATH=/var/www/iqoptiontraderbot/venv/bin"
ExecStart=/var/www/iqoptiontraderbot/venv/bin/gunicorn \
    --config /var/www/iqoptiontraderbot/gunicorn.conf.py \
    wsgi:app

Restart=always
RestartSec=3

# Segurança
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/www/iqoptiontraderbot/logs /var/www/iqoptiontraderbot/sinais.txt

[Install]
WantedBy=multi-user.target

